# system/mab-api.service
[Unit]
Description=MAB API Service
After=network.target

[Service]
Type=simple
User=william
WorkingDirectory=/home/william/personal/mab
ExecStart=/home/william/.bun/bin/bun api-server.js
Environment=PORT=8404
Environment=NODE_ENV=development
StandardInput=null
StandardOutput=journal
StandardError=journal
Restart=no

[Install]
WantedBy=multi-user.target
